# examples/multi-quadlets/hello.container
# hello.container
[Container]
Image=nginx.image
PublishPort=8888:80

[Service]
Restart=always
